# etc/systemd/system/emp-api.service
[Unit]
Description=United mobile app API - Patient Management FastAPI Service
After=network.target
Wants=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/emp_api/api_qms_emp
Environment=PATH=/usr/bin:/bin
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 7072
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=emp-api

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/emp_api/api_qms_emp
ProtectHome=yes

[Install]
WantedBy=multi-user.target